# config/supervisord.conf
[supervisord]
nodaemon=true

[program:pmtahttpd]
command=pmtahttpd

[program:pmtad]
command=pmtad
